# systemd/led-status.service
[Unit]
Description=Pi Assistant LED Status Ring
After=network.target orchestrator.service
Wants=orchestrator.service

[Service]
Type=simple
WorkingDirectory=/home/dev/project_root
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/dev/project_root/.env
ExecStart=/home/dev/project_root/.venvs/visn/bin/python -m src.piled.led_ring_service
Restart=always
RestartSec=3
StandardOutput=append:/home/dev/project_root/logs/led_status.log
StandardError=append:/home/dev/project_root/logs/led_status.service.log

[Install]
WantedBy=multi-user.target
