# systemd/led-ring.service
[Unit]
Description=Smart Car LED Ring Status (NeoPixel)
After=orchestrator.service
Requires=orchestrator.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
WorkingDirectory=/home/dev/smart_car
Environment=PROJECT_ROOT=/home/dev/smart_car
Environment=PYTHONPATH=/home/dev/smart_car
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/dev/smart_car/.env
ExecStart=/home/dev/smart_car/.venvs/visn-py313/bin/python -m src.piled.led_ring_service --config config/system.yaml
Restart=on-failure
RestartSec=3
StandardOutput=append:/home/dev/smart_car/logs/led-ring.log
StandardError=append:/home/dev/smart_car/logs/led-ring.log

[Install]
WantedBy=multi-user.target
