# systemd/stt-wrapper.service
 [Unit]
Description=STT Wrapper service (AudioManager-based STT)
After=network.target audio-manager.service

[Service]
Type=simple
WorkingDirectory=/home/dev/project_root
ExecStart=/home/dev/project_root/.venvs/stte/bin/python -m src.stt.stt_wrapper_runner --config config/system.yaml
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
