# systemd/tts.service
[Unit]
Description=Smart Car TTS (Azure)
After=orchestrator.service
Requires=orchestrator.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=dev
Group=dev
WorkingDirectory=/home/dev/smart_car
Environment=PROJECT_ROOT=/home/dev/smart_car
Environment=PYTHONPATH=/home/dev/smart_car
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/dev/smart_car/.env
ExecStart=/home/dev/smart_car/.venvs/ttse/bin/python -m src.tts.azure_tts_runner
Restart=on-failure
RestartSec=3
StandardOutput=append:/home/dev/smart_car/logs/tts.log
StandardError=append:/home/dev/smart_car/logs/tts.log

[Install]
WantedBy=multi-user.target
