# systemd/wakeword.service
[Unit]
Description=Pi Assistant Wakeword Service (Porcupine dsnoop)
After=network.target sound.target

[Service]
Type=simple
WorkingDirectory=/home/dev/project_root
Environment=PV_ACCESS_KEY=${PV_ACCESS_KEY}
Environment=PYTHONPATH=/home/dev/project_root
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/dev/project_root/.env
ExecStart=/home/dev/project_root/.venvs/stte/bin/python -m src.wake.service --config config/system.yaml
Restart=always
RestartSec=3
TimeoutStartSec=30
User=dev
Group=dev
StandardOutput=append:/home/dev/project_root/logs/wakeword.log
StandardError=append:/home/dev/project_root/logs/wakeword.service.log

[Install]
WantedBy=multi-user.target
